# cron/crontab
# FinOps reminder cron — runs daily at 9:00 AM (container local time)
# Adjust timezone via TZ env var in docker-compose
0 9 * * * wget -q -O- http://finops:3000/api/remind > /dev/null 2>&1
